# crontab
# Pump Detection System Crontab
# Install with: crontab /home/elcrypto/pump_detector/crontab

# Daily calibration at 2 AM
0 2 * * * cd /home/elcrypto/pump_detector && /usr/bin/python3 scripts/calibrate_scoring.py >> logs/calibration.log 2>&1

# Weekly report generation (Sunday at 3 AM)
0 3 * * 0 cd /home/elcrypto/pump_detector && /usr/bin/python3 scripts/generate_reports.py >> logs/reports.log 2>&1

# Hourly OI data tracking
0 * * * * cd /home/elcrypto/pump_detector && psql -d fox_crypto_new -c "SELECT pump.track_oi_changes();" >> logs/oi_tracking.log 2>&1

# Daily cleanup of old signals (keep 30 days)
0 4 * * * psql -d fox_crypto_new -c "DELETE FROM pump.signal_tracking WHERE check_timestamp < NOW() - INTERVAL '30 days';" >> /home/elcrypto/pump_detector/logs/cleanup.log 2>&1

# Check daemon health every 5 minutes
*/5 * * * * /home/elcrypto/pump_detector/scripts/check_daemons.sh

# Backup database daily at 5 AM
0 5 * * * pg_dump -d fox_crypto_new -n pump -f /home/elcrypto/pump_detector/backups/pump_$(date +\%Y\%m\%d).sql 2>> /home/elcrypto/pump_detector/logs/backup.log
